# forensic-dashboard.service
[Unit]
Description=Forensic Suite Dashboard - Web Interface
After=network.target

[Service]
Type=simple
User=hackerguarro
Group=hackerguarro
WorkingDirectory=/usr/local/Forensic-Suite-dashboard
Environment="PATH=/usr/local/Forensic-Suite-dashboard/.venv/bin"
ExecStart=/usr/local/Forensic-Suite-dashboard/.venv/bin/python app.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
